# cron/crontab
# Cron job to run full backup script every week on Sunday at 3:00 AM
0 3 * * 0 /usr/local/bin/node /app/scripts/fullBackup.js >> /var/log/cron.log 2>&1